# etc/systemd/system/weather.service
[Unit]
Description=Update /etc/motd with weather from wttr.in
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/default/weather
ExecStart=/usr/local/sbin/weather.sh --today
RemainAfterExit=yes
# (Bonus: activează logare în fișier – decomentează liniile de mai jos)
# StandardOutput=append:/var/log/weather.log
# StandardError=append:/var/log/weather.log

[Install]
WantedBy=multi-user.target
